# scripts/init.d/vpn-mode
#!/bin/sh /etc/rc.common
# При загрузке восстанавливает режим VPN.
# На Beryl AX — читает положение GPIO-слайдера (vpn-mode detect).
# На Cudy TR3000 и других — применяет последний сохранённый режим из /etc/vpn-mode.state.

START=99
USE_PROCD=1

boot() { start; }

start_service() {
    # Ждём podkop, потом синхронизируем режим
    ( sleep 5; /usr/bin/vpn-mode detect ) &
}
